# software/cloud-dashboard/supervisord.conf
[supervisord]
nodaemon=true

[program:mongod]
command=/usr/bin/mongod

[program:node]
command=bash -c 'sleep 15 && node /src/server.js'
