# software/cloud-dashboard/supervisord.conf
[supervisord]
nodaemon=true

[program:mongod]
command=/usr/bin/mongod

[program:node]
command=node /src/ADPL/software/cloud-dashboard/node-server/cloud-server.js
