# software/logging/supervisord.conf
[supervisord]
nodaemon=true

[program:mongod]
command=/usr/bin/mongod

[program:node]
command=node /src/server.js
